# vm/init-vm.service
[Unit]
Description=VM init service

[Service]
ExecStart=/usr/bin/_entry_vm_init

[Install]
WantedBy=multi-user.target